# qualifiers/deploy/ansible/roles/pg-gcs-backup/files/backup-db.timer
[Unit]
Description=Timer to run backup-db.service

[Timer]
Unit=backup-db.service
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
