# qualifiers/deploy/ansible/roles/pg-gcs-backup/files/backup-db.service
[Unit]
Description=Backup Postgres database and upload to GCS

[Service]
User=postgres
Group=postgres
ExecStart=/opt/backup/backup-db.sh

[Install]
WantedBy=multi-user.target
