# qualifiers/deploy/ansible/roles/score-reporter/files/score-reporter.timer
[Unit]
Description=Timer to run score-reporter.service

[Timer]
Unit=score-reporter.service
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
